# REE-DEB/DEBIAN/preinst
#!/bin/bash -u

# Activation du wifi
sudo rfkill unblock wifi

sudo mkdir -p /usr/local/bin &> /dev/null
sudo mkdir -p /usr/local/default &> /dev/null
sudo mkdir -p /usr/local/security &> /dev/null
sudo mkdir -p /usr/local/vars &> /dev/null
sudo mv ../usr/local/* /usr/local/
chmod u+x /usr/local/bin/*


# Installation des paquets
apt install iptables -y
apt install iptables-persistent -y
apt install dnsmasq -y
apt install hostapd -y
apt install dhcpcd -y
apt install espeak -y


# Configuration wpa_supplicant pour wlan1 (= cherche à se connecter au réseau wifi REE)
cat /usr/local/default/wpa_supplicant-wlan1.conf > /etc/wpa_supplicant/wpa_supplicant-wlan1.conf

systemctl stop wpa_supplicant.service
systemctl disable wpa_supplicant.service

systemctl start wpa_supplicant@wlan1.service
systemctl enable wpa_supplicant@wlan1.service
systemctl status wpa_supplicant@wlan1.service

dhclient

# Configuration hostapd 
cat /usr/local/default/hostapd.conf > /etc/hostapd/hostapd.conf

systemctl unmask hostapd
systemctl enable hostapd
systemctl restart hostapd.service






# sudo useradd -m rasp
# sudo groupadd security
# sudo usermod -aG security rasp

# sudo chown -R root:security /srv/security/*
