# REE-DEB/DEBIAN/preinst
#/bin/bash -u

# Installation des paquets
apt install dnsmasq -y
apt install hostapd -y
apt install iptables -y
apt install iptables-persistent -y

# Configuration wpa_supplicant
echo 'network={
    ssid="REE"
    psk="REE123456" #A modifier avec le mot de passe du Wi-Fi REE
}' >> /etc/wpa_supplicant/wpa_supplicant-wlan1.conf

systemctl stop wpa_supplicant.service
systemctl disable wpa_supplicant.service

systemctl start wpa_supplicant@wlan1.service
systemctl enable wpa_supplicant@wlan1.service
systemctl status wpa_supplicant@wlan1.service

dhclient

# Configuration inteface wlan0 
#/!\ automatisé fonction de Pi paire ou impaire
echo -n 'interface wlan0
static ip_address=10.0.1.1/30
#IP DE la box
static routers=192.168.1.254
static domain_name_servers=8.8.8.8
' >> /etc/dhcpcd.conf

systemctl restart dhcpcd.service


# Configuration DNSmasq
#/!\ automatisé fonction de Pi paire ou impaire (sois 10.0.1.2 et 10.0.1.3 soit 10.0.2.2 et 10.0.2.3)


echo 'interface=wlan0
dhcp-range=10.0.1.2,10.0.1.3,255.255.255.252,24h
server=8.8.8.8
' >> /etc/dnsmasq.conf

systemctl restart dnsmasq.service

# Configuration hostapd 
echo 'interface=wlan0
ssid=REE
hw_mode=g
channel=1
macaddr_acl=0
auth_algs=1
#ignore_broadcast_ssid=0
wpa=2
wpa_passphrase=REE123456789
wpa_key_mgmt=WPA-PSK
wpa_pairwise=TKIP
rsn_pairwise=CCMP
country_code=FR
#ctrl_interface=/var/run/hostapd
#ctrl_interface_group=0
#daemonize /usr/local/bin/manageroutes
' > /etc/hostapd/hostapd.conf

sudo systemctl unmask hostapd
sudo systemctl enable hostapd
systemctl restart hostapd.service

# /!\ Voir si on ne peut pas empêcher hostapd d'ecrire ces propre routes
# Que pour les PI wireless -> wireless
#/!\ automatisé fonction de Pi paire ou impaire
ip route del default via 10.0.3.254 dev wlan0 src 11.0.0.1 metric 3003
ip route del 10.0.3.254 dev wlan0 scope link src 11.0.0.1 metric 3003

#10.0.1.1 #non-Wlan0 IP
ip route add -net 10.0.1.1 netmask 255.255.255.252 dev wlan0 


# Configuration iptables

# echo '*filter 
# :INPUT ACCEPT [0:0] 
# :FORWARD ACCEPT [0:0] 
# :OUTPUT ACCEPT [0:0] 
# -A FORWARD -i wlan0 -j ACCEPT 
# -A FORWARD -o wlan0 -j ACCEPT 
# COMMIT 
# *nat 
# :PREROUTING ACCEPT [0:0] 
# :INPUT ACCEPT [0:0] 
# :OUTPUT ACCEPT [4:324] 
# :POSTROUTING ACCEPT [0:0] 
# -A POSTROUTING -o wlan1 -j MASQUERADE 
# COMMIT' > /etc/iptables/rules.v4

# iptables -F
# iptables -t nat -F
# iptables -P INPUT ACCEPT
# iptables -P FORWARD ACCEPT
# iptables -P OUTPUT ACCEPT
# iptables -A FORWARD -i wlan0 -o eth0 -j ACCEPT
# iptables -A FORWARD -i eth0 -o wlan0 -j ACCEPT
# iptables -t nat -A POSTROUTING -o eth0 -j MASQUERADE
# sh -c "iptables-save > /etc/iptables/rules.v4" 
